# cron/patch-system.cron
# Weekly system patching
# Runs every Friday at 03:15 AM

15 3 * * 5 root  /usr/local/bin/patch-system.sh >> /var/log/patch-system.log 2>&1
